# fastfetch_conf/systemd_jobs/notify-endeavour.service
[Unit]
Description=Envía notificación para actualizar Endeavour

[Service]
Type=oneshot
ExecStart=/home/patricioeb/.local/bin/notify-endeavour.sh
