# systemd_jobs/notify-hour.timer
[Unit]
Description=Timer para notificación horaria

[Timer]
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target
