# systemd_jobs/notify-hour.service
[Unit]
Description=Envía notificación con la hora actual

[Service]
Type=oneshot
ExecStart=/home/patricioeb/.local/bin/notify-hour.sh
